# server/systemd/user/sync-arch.timer
[Unit]
Description=Sync package repositories every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
