# server/systemd/system/libvirtd.service.d/override.conf
# Otherwise "Cannot get interface MTU on 'br0': No such device" will prevent VMs
# autostarting
[Unit]
Wants=network-online.target
After=network-online.target
