# server/systemd/system/squid.service.d/override.conf
# The default service file adds -s ("enable logging to syslog") and fills the
# journal with noise
[Service]
ExecStart=
ExecStart=/usr/sbin/squid --foreground -YC
